# .config/systemd/user/dunst.service
[Unit]
Description=dunst notification service

[Service]
Type=simple
ExecStart=/usr/bin/dunst

[Install]
WantedBy=default.target
