# .config/systemd/user/imwheel.service
[Unit]
Description=imwheel daemon

[Service]
Type=simple
ExecStart=/usr/bin/imwheel -b 45

[Install]
WantedBy=default.target
